# ops/thrushcombe.timer
[Unit]
Description=Hourly Thrushcombe tick (companion mode; the day has phases)

[Timer]
# Every hour at :07 — keeps the town caught up to the current phase of the day, so
# morning, market, afternoon-calls and evening-pub beats land as the day passes.
# Persistent fires it on next boot if the machine was off.
OnCalendar=*-*-* *:07:00
Persistent=true

[Install]
WantedBy=timers.target
